# system/custom-skeleton/ubuntu-base-16.04.2-arm64/etc/init/hwclock.sh.conf
description "Mask /etc/init.d/hwclock.sh, this is done by /etc/init/hwclock.conf"

start on starting hwclock
