# system/custom-skeleton/ubuntu-base-16.04.2-arm64/etc/init/procps.conf
# procps - set sysctls from /etc/sysctl.conf
#
# signal startpar that procps has run

start on started procps-instance
