# board/nxp/ls1021atsn/rootfs_overlay/usr/lib/systemd/system/phc2sys.service
[Unit]
Description=Precision Time Protocol system clock synchronization
After=ptp4l.service

[Service]
ExecStart=/usr/sbin/phc2sys -a -rr -f /etc/linuxptp-system-clock.cfg
Restart=always

[Install]
WantedBy=multi-user.target
