# board/nxp/common/systemd_init_overlay/usr/lib/systemd/system/serial-getty@.service.d/default-TERM.conf
[Service]
Environment=TERM=screen-256color
