# system/custom-skeleton/ubuntu-base-16.04.2-arm64/lib/systemd/system/sigpwr-container-shutdown.service
[Unit]
Description=Shut down containers on SIGPWR
ConditionVirtualization=container

[Service]
Type=oneshot
ExecStart=/bin/systemctl --no-block poweroff
