# system/custom-skeleton/ubuntu-base-16.04.2-arm64/lib/systemd/system/systemd-networkd-resolvconf-update.path
[Unit]
Description=Trigger resolvconf update for networkd DNS
ConditionPathIsSymbolicLink=/etc/resolv.conf

[Path]
PathChanged=/run/systemd/netif/state
