# system/custom-skeleton/ubuntu-base-16.04.2-arm64/etc/cron.weekly/fstrim
#!/bin/sh
# trim all mounted file systems which support it
/sbin/fstrim --all || true
