# system/custom-skeleton/ubuntu-base-16.04.2-arm64/var/lib/dpkg/info/init-system-helpers.postinst
#!/bin/sh

set -e

# See http://bugs.debian.org/722524
if dpkg --compare-versions "$2" lt "1.9"; then
    if [ -d /var/lib/systemd/deb-systemd-helper-enabled ]; then
        find /var/lib/systemd/deb-systemd-helper-enabled -type f -exec chmod 644 '{}' \;
    fi
fi



exit 0
