# system/custom-skeleton/ubuntu-base-16.04.2-arm64/var/lib/dpkg/info/sysv-rc.list
/.
/etc
/etc/rc0.d
/etc/rc0.d/README
/etc/init.d
/etc/init.d/rcS
/etc/init.d/README
/etc/init.d/rc
/etc/rc3.d
/etc/rc3.d/README
/etc/rc1.d
/etc/rc1.d/README
/etc/rcS.d
/etc/rcS.d/README
/etc/rc2.d
/etc/rc2.d/README
/etc/rc6.d
/etc/rc6.d/README
/etc/rc5.d
/etc/rc5.d/README
/etc/rc4.d
/etc/rc4.d/README
/var
/var/lib
/var/lib/insserv
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/sysv-rc
/usr/share/doc
/usr/share/doc/sysv-rc
/usr/share/doc/sysv-rc/copyright
/usr/share/doc/sysv-rc/README.Debian
/usr/share/doc/sysv-rc/NEWS.Debian.gz
/usr/share/doc/sysv-rc/README.runlevels.gz
/usr/share/sysv-rc
/usr/share/sysv-rc/saveconfig
/usr/share/doc/sysv-rc/changelog.Debian.gz
