# system/custom-skeleton/ubuntu-base-16.04.2-arm64/var/lib/dpkg/info/util-linux.list
/.
/sbin
/sbin/switch_root
/sbin/ctrlaltdel
/sbin/mkfs.bfs
/sbin/mkfs
/sbin/swaplabel
/sbin/chcpu
/sbin/fsck.cramfs
/sbin/fsck
/sbin/zramctl
/sbin/fdisk
/sbin/pivot_root
/sbin/fsfreeze
/sbin/blkdiscard
/sbin/raw
/sbin/runuser
/sbin/blkid
/sbin/mkswap
/sbin/cfdisk
/sbin/wipefs
/sbin/sulogin
/sbin/findfs
/sbin/mkfs.minix
/sbin/blockdev
/sbin/mkfs.cramfs
/sbin/sfdisk
/sbin/fstrim
/sbin/hwclock
/sbin/fsck.minix
/sbin/agetty
/sbin/isosize
/usr
/usr/sbin
/usr/sbin/fdformat
/usr/sbin/readprofile
/usr/sbin/tunelp
/usr/sbin/ldattach
/usr/sbin/rtcwake
/usr/bin
/usr/bin/rename.ul
/usr/bin/flock
/usr/bin/mesg
/usr/bin/rev
/usr/bin/utmpdump
/usr/bin/mcookie
/usr/bin/line
/usr/bin/ipcs
/usr/bin/ipcmk
/usr/bin/ipcrm
/usr/bin/resizepart
/usr/bin/partx
/usr/bin/prlimit
/usr/bin/taskset
/usr/bin/getopt
/usr/bin/addpart
/usr/bin/whereis
/usr/bin/lscpu
/usr/bin/pg
/usr/bin/unshare
/usr/bin/delpart
/usr/bin/nsenter
/usr/bin/chrt
/usr/bin/setsid
/usr/bin/ionice
/usr/bin/lslogins
/usr/bin/last
/usr/bin/setarch
/usr/bin/lsipc
/usr/bin/fallocate
/usr/bin/lslocks
/usr/bin/setterm
/usr/bin/namei
/usr/share
/usr/share/doc
/usr/share/doc/util-linux
/usr/share/doc/util-linux/howto-compilation.txt
/usr/share/doc/util-linux/NEWS.Debian.gz
/usr/share/doc/util-linux/poeigl.txt.gz
/usr/share/doc/util-linux/howto-man-page.txt.gz
/usr/share/doc/util-linux/PAM-configuration.txt
/usr/share/doc/util-linux/AUTHORS.gz
/usr/share/doc/util-linux/mount.txt
/usr/share/doc/util-linux/release-schedule.txt
/usr/share/doc/util-linux/deprecated.txt
/usr/share/doc/util-linux/blkid.txt
/usr/share/doc/util-linux/source-code-management.txt
/usr/share/doc/util-linux/col.txt
/usr/share/doc/util-linux/howto-usage-function.txt.gz
/usr/share/doc/util-linux/examples
/usr/share/doc/util-linux/examples/udev-raw.rules
/usr/share/doc/util-linux/examples/shells
/usr/share/doc/util-linux/examples/getopt-parse.bash
/usr/share/doc/util-linux/examples/fstab
/usr/share/doc/util-linux/examples/fstrim.timer
/usr/share/doc/util-linux/examples/filesystems
/usr/share/doc/util-linux/examples/getopt-parse.tcsh
/usr/share/doc/util-linux/examples/motd
/usr/share/doc/util-linux/examples/fstrim.service
/usr/share/doc/util-linux/examples/securetty
/usr/share/doc/util-linux/examples/fstab.example2
/usr/share/doc/util-linux/releases
/usr/share/doc/util-linux/releases/v2.20-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.24-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.18-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.13-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.25-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.17-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.26-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.19-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.16-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.27-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.14-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.23-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.27.1-ReleaseNotes
/usr/share/doc/util-linux/releases/v2.21-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.15-ReleaseNotes.gz
/usr/share/doc/util-linux/releases/v2.22-ReleaseNotes.gz
/usr/share/doc/util-linux/howto-pull-request.txt.gz
/usr/share/doc/util-linux/getopt_changelog.txt
/usr/share/doc/util-linux/howto-contribute.txt.gz
/usr/share/doc/util-linux/howto-debug.txt
/usr/share/doc/util-linux/00-about-docs.txt
/usr/share/doc/util-linux/cal.txt
/usr/share/doc/util-linux/howto-tests.txt
/usr/share/doc/util-linux/hwclock.txt
/usr/share/doc/util-linux/copyright
/usr/share/doc/util-linux/modems-with-agetty.txt
/usr/share/doc/util-linux/howto-build-sys.txt
/usr/share/doc/util-linux/pg.txt
/usr/share/doc/util-linux/getopt.txt
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/flock
/usr/share/bash-completion/completions/mesg
/usr/share/bash-completion/completions/rev
/usr/share/bash-completion/completions/utmpdump
/usr/share/bash-completion/completions/mcookie
/usr/share/bash-completion/completions/ctrlaltdel
/usr/share/bash-completion/completions/mkfs.bfs
/usr/share/bash-completion/completions/mkfs
/usr/share/bash-completion/completions/swaplabel
/usr/share/bash-completion/completions/chcpu
/usr/share/bash-completion/completions/ipcs
/usr/share/bash-completion/completions/fsck.cramfs
/usr/share/bash-completion/completions/fsck
/usr/share/bash-completion/completions/zramctl
/usr/share/bash-completion/completions/more
/usr/share/bash-completion/completions/wdctl
/usr/share/bash-completion/completions/fdisk
/usr/share/bash-completion/completions/pivot_root
/usr/share/bash-completion/completions/fdformat
/usr/share/bash-completion/completions/ipcrm
/usr/share/bash-completion/completions/resizepart
/usr/share/bash-completion/completions/partx
/usr/share/bash-completion/completions/prlimit
/usr/share/bash-completion/completions/taskset
/usr/share/bash-completion/completions/fsfreeze
/usr/share/bash-completion/completions/getopt
/usr/share/bash-completion/completions/blkdiscard
/usr/share/bash-completion/completions/readprofile
/usr/share/bash-completion/completions/raw
/usr/share/bash-completion/completions/addpart
/usr/share/bash-completion/completions/whereis
/usr/share/bash-completion/completions/lscpu
/usr/share/bash-completion/completions/mountpoint
/usr/share/bash-completion/completions/runuser
/usr/share/bash-completion/completions/blkid
/usr/share/bash-completion/completions/tailf
/usr/share/bash-completion/completions/mkswap
/usr/share/bash-completion/completions/pg
/usr/share/bash-completion/completions/unshare
/usr/share/bash-completion/completions/cfdisk
/usr/share/bash-completion/completions/delpart
/usr/share/bash-completion/completions/wipefs
/usr/share/bash-completion/completions/nsenter
/usr/share/bash-completion/completions/mkfs.minix
/usr/share/bash-completion/completions/chrt
/usr/share/bash-completion/completions/setsid
/usr/share/bash-completion/completions/ionice
/usr/share/bash-completion/completions/tunelp
/usr/share/bash-completion/completions/blockdev
/usr/share/bash-completion/completions/mkfs.cramfs
/usr/share/bash-completion/completions/sfdisk
/usr/share/bash-completion/completions/ldattach
/usr/share/bash-completion/completions/fstrim
/usr/share/bash-completion/completions/last
/usr/share/bash-completion/completions/rtcwake
/usr/share/bash-completion/completions/hwclock
/usr/share/bash-completion/completions/setarch
/usr/share/bash-completion/completions/lsipc
/usr/share/bash-completion/completions/fallocate
/usr/share/bash-completion/completions/dmesg
/usr/share/bash-completion/completions/lslocks
/usr/share/bash-completion/completions/fsck.minix
/usr/share/bash-completion/completions/setterm
/usr/share/bash-completion/completions/lsblk
/usr/share/bash-completion/completions/isosize
/usr/share/bash-completion/completions/namei
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/fsfreeze.8.gz
/usr/share/man/man8/setarch.8.gz
/usr/share/man/man8/tunelp.8.gz
/usr/share/man/man8/lslocks.8.gz
/usr/share/man/man8/rtcwake.8.gz
/usr/share/man/man8/isosize.8.gz
/usr/share/man/man8/zramctl.8.gz
/usr/share/man/man8/fdformat.8.gz
/usr/share/man/man8/pivot_root.8.gz
/usr/share/man/man8/blkid.8.gz
/usr/share/man/man8/mkfs.minix.8.gz
/usr/share/man/man8/cfdisk.8.gz
/usr/share/man/man8/mkswap.8.gz
/usr/share/man/man8/wdctl.8.gz
/usr/share/man/man8/agetty.8.gz
/usr/share/man/man8/hwclock.8.gz
/usr/share/man/man8/sulogin.8.gz
/usr/share/man/man8/fsck.8.gz
/usr/share/man/man8/ctrlaltdel.8.gz
/usr/share/man/man8/mkfs.bfs.8.gz
/usr/share/man/man8/blkdiscard.8.gz
/usr/share/man/man8/wipefs.8.gz
/usr/share/man/man8/swaplabel.8.gz
/usr/share/man/man8/sfdisk.8.gz
/usr/share/man/man8/delpart.8.gz
/usr/share/man/man8/ldattach.8.gz
/usr/share/man/man8/fsck.minix.8.gz
/usr/share/man/man8/findfs.8.gz
/usr/share/man/man8/fdisk.8.gz
/usr/share/man/man8/blockdev.8.gz
/usr/share/man/man8/resizepart.8.gz
/usr/share/man/man8/lsblk.8.gz
/usr/share/man/man8/mkfs.cramfs.8.gz
/usr/share/man/man8/fsck.cramfs.8.gz
/usr/share/man/man8/addpart.8.gz
/usr/share/man/man8/fstrim.8.gz
/usr/share/man/man8/switch_root.8.gz
/usr/share/man/man8/readprofile.8.gz
/usr/share/man/man8/mkfs.8.gz
/usr/share/man/man8/chcpu.8.gz
/usr/share/man/man8/partx.8.gz
/usr/share/man/man8/raw.8.gz
/usr/share/man/man5
/usr/share/man/man5/hwclock.5.gz
/usr/share/man/man5/terminal-colors.d.5.gz
/usr/share/man/man1
/usr/share/man/man1/dmesg.1.gz
/usr/share/man/man1/whereis.1.gz
/usr/share/man/man1/setsid.1.gz
/usr/share/man/man1/chrt.1.gz
/usr/share/man/man1/rename.ul.1.gz
/usr/share/man/man1/flock.1.gz
/usr/share/man/man1/line.1.gz
/usr/share/man/man1/lsipc.1.gz
/usr/share/man/man1/rev.1.gz
/usr/share/man/man1/ipcs.1.gz
/usr/share/man/man1/nsenter.1.gz
/usr/share/man/man1/mcookie.1.gz
/usr/share/man/man1/getopt.1.gz
/usr/share/man/man1/ipcrm.1.gz
/usr/share/man/man1/utmpdump.1.gz
/usr/share/man/man1/fallocate.1.gz
/usr/share/man/man1/lslogins.1.gz
/usr/share/man/man1/runuser.1.gz
/usr/share/man/man1/ipcmk.1.gz
/usr/share/man/man1/pg.1.gz
/usr/share/man/man1/prlimit.1.gz
/usr/share/man/man1/unshare.1.gz
/usr/share/man/man1/mountpoint.1.gz
/usr/share/man/man1/tailf.1.gz
/usr/share/man/man1/last.1.gz
/usr/share/man/man1/taskset.1.gz
/usr/share/man/man1/namei.1.gz
/usr/share/man/man1/setterm.1.gz
/usr/share/man/man1/lscpu.1.gz
/usr/share/man/man1/ionice.1.gz
/usr/share/man/man1/more.1.gz
/usr/share/man/man1/mesg.1.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/util-linux
/usr/lib
/usr/lib/mime
/usr/lib/mime/packages
/usr/lib/mime/packages/util-linux
/bin
/bin/more
/bin/wdctl
/bin/mountpoint
/bin/tailf
/bin/dmesg
/bin/lsblk
/etc
/etc/init.d
/etc/init.d/hwclock.sh
/etc/init
/etc/init/hwclock-save.conf
/etc/init/hwclock.conf
/etc/init/hwclock.sh.conf
/etc/pam.d
/etc/pam.d/runuser-l
/etc/pam.d/runuser
/etc/cron.weekly
/etc/cron.weekly/fstrim
/lib
/lib/udev
/lib/udev/hwclock-set
/sbin/getty
/usr/bin/linux64
/usr/bin/lastb
/usr/bin/linux32
/usr/share/doc/util-linux/changelog.Debian.gz
/usr/share/man/man8/getty.8.gz
/usr/share/man/man1/linux64.1.gz
/usr/share/man/man1/linux32.1.gz
/usr/share/man/man1/lastb.1.gz
